# deploy/cortex-server.service
[Unit]
Description=Cortex Web UI Server
After=network-online.target cortex-daemon.service
Wants=network-online.target cortex-daemon.service

[Service]
Type=simple
ExecStart=/usr/local/bin/cortex server start --port 3000 --host 127.0.0.1
Restart=always
RestartSec=5
Environment="HOME=%h"

[Install]
WantedBy=default.target
